# .config/rolr/roles/aarch64/ubuntu/desktop-kiosk/udev.rules
# managed by rolr
ACTION=="add", SUBSYSTEM=="backlight", KERNEL=="10-0045", RUN+="/bin/sh -c '/bin/chmod 666 /sys/class/backlight/%k/brightness'"
